# batchsvc/crontab
* * * * * root /opt/batch_service/runbatch.sh >> /var/log/cron.log
